# deploy/antiddosd.service
[Unit]
Description=XDP Anti-DDoS mitigator
After=network-online.target
Wants=network-online.target
Documentation=https://github.com/anti-ddos/antiddosd

[Service]
Type=simple
ExecStartPre=/usr/local/sbin/antiddos-tuning %i %I
ExecStart=/usr/local/bin/antiddosd \
    --config /etc/antiddosd/config.yaml \
    --log-json
Restart=on-failure
RestartSec=2s
# Required for loading BPF programs and attaching XDP.
AmbientCapabilities=CAP_NET_ADMIN CAP_BPF CAP_PERFMON CAP_SYS_ADMIN
CapabilityBoundingSet=CAP_NET_ADMIN CAP_BPF CAP_PERFMON CAP_SYS_ADMIN
NoNewPrivileges=true
LimitMEMLOCK=infinity
LimitNOFILE=1048576
# Isolate to non-hot cores; datapath runs in softirq on the NIC IRQ cores.
CPUAffinity=0
OOMScoreAdjust=-900
TasksMax=512
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
